# conf/rc.local
#!/bin/sh -e
cd /home/ubuntu && sudo -u ubuntu jupyter notebook --no-browser --ip=0.0.0.0 --NotebookApp.token="" &
#rstudio-server start
exit 0